# .crontab
# check doi requests
* */12 * * * [PROJECT_PATH]/app/console vipa:article:doi:requests:check
